# x2go/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:bash]
command=/bin/bash /run.sh
startsecs = 0
autorestart = false
startretries = 1

[program:ssh]
command=/usr/sbin/sshd -D
